# etc/cron.d/hsec-dosportal-graylog-web
#
# Managed by Ansible
#
# get IPs from graylog & insert them into dosportal
*/5 * * * * root /localservices/sbin/lockrun /var/run/dosportal-graylog-web.lock /localservices/sbin/dosportal-graylog-web.pl > /dev/null 2>&1

